# tools/centos/etc/systemd/system/server.service
[Unit]
Description=Server
After=syslog.target network.target

[Service]
Type=simple
User=root
ExecStart=/data/wwwroot/server/main
Restart=on-failure

[Install]
WantedBy=multi-user.target
